# board/licheervnano/rootfs_overlay/etc/inittab
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -o remount,rw /
::sysinit:/bin/mkdir -p /dev/pts /dev/shm
::sysinit:/bin/mount -a
::sysinit:/bin/hostname -F /etc/hostname
::sysinit:/etc/init.d/rcS

::respawn:/sbin/getty -L ttyS0 115200 vt100
::respawn:/sbin/getty -L ttyGS0 115200 vt100

::shutdown:/etc/init.d/rcK
::ctrlaltdel:/sbin/reboot
